# files/crontab
# Pull and Replace pihole image every week
@weekly compose-update.sh
